# SOURCES/app_b-1/app_b.init
#!/bin/bash
#
# app_b
#
# chkconfig: 345 70 30
# description: density hw app_b
# processname: app_b

# Source function library.
#. /etc/init.d/functions

RETVAL=0
NAME="app_b"
PIDFILE=/var/run/$NAME
OS_VER=`grep -o '[[:digit:]]' /etc/redhat-release | head -1`

start() {
        echo -n "Starting $NAME: "
        if [ ! -f /opt/density/database.db ]; then
	  sqlite3 /opt/density/database.db < /opt/density/schema.sql
	fi
	if [ $OS_VER -lt 8 ]; then
	  python /opt/density/app_b.py
	else
	  python3 /opt/density/app_b.py
	fi
        RETVAL=$?
        [ $RETVAL -eq 0 ] && touch $PIDFILE
        echo
        return $RETVAL
}

stop() {
        echo -n "Shutting down $NAME: "
        killall python
        rm -f $PIDFILE
        RETVAL=$?
        [ $RETVAL -eq 0 ] && rm -f $PIDFILE
        echo
        return $RETVAL
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo "Usage: $NAME {start|stop|restart}"
        exit 1
        ;;
esac
exit $RETVAL
